# ahrimod.service
[Unit]
Description=Ahrimod — Telegram moderation bot for AhriVPN
Documentation=https://github.com/kifchandr/AhriMOD
After=network-online.target
Wants=network-online.target
# Лимит на частоту перезапусков: не больше 5 крашей за 5 минут,
# иначе systemd пометит сервис как failed и перестанет рестартить.
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple

# Запускается от root (т.к. лежит в /root/AhriMOD)
User=root
Group=root

# Где лежит проект
WorkingDirectory=/root/AhriMOD

# Запуск через venv
ExecStart=/root/AhriMOD/venv/bin/python /root/AhriMOD/main.py

# Перезапуск при падениях
Restart=always
RestartSec=10

# Окружение. FORCE_COLOR=1 включает цветной rich-вывод в journald,
# который иначе считается не-TTY и режет ANSI-коды.
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=FORCE_COLOR=1
Environment=TERM=xterm-256color

# Логирование в journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ahrimod

# Таймаут штатной остановки
TimeoutStopSec=15
KillMode=mixed
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
